# rootfs/etc/profile.d/env_pre.sh
echo "customer env at: etc/profile.d"
echo "source env from /system/env/env.sh"
source /system/env/env.sh
