# rootfs/etc/profile.d/env_pre.sh
echo "customer env at: etc/profile.d"
source /system/env_export.sh
echo "start Android shell"
sh
